# DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop opensensecam.service
systemctl disable opensensecam.service

exit 0